# crontab
*/2 * * * * root /usr/local/bin/python /home/mlflowserver/mlflow/test.py >> /var/log/cron.log 2>&1
#
